# services/reboot_hung_systems/reboot_hung_systems.service
[Unit]
Description=Control managed switches to reboot hung systems

[Service]
Type=simple
ExecStart=/var/cache/mesa_jenkins/repos/mesa_ci/services/reboot_hung_systems/reboot_hung_systems.service
Restart=on-failure

[Install]
WantedBy=multi-user.target
